# infrastructure/tensorzero/health/tensorzero-health.service
[Unit]
Description=TensorZero health check and auto-recovery
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tensorzero-health-check.sh
User=root
